# flex-receiver-flexdec.service
[Unit]
Description=FLEX paging receiver (RSPdx -> in-memory flexdec, ALN pages)
After=network.target sdrplay.service
Wants=sdrplay.service

[Service]
Type=simple
User=swannman
ExecStart=/usr/bin/python3 /usr/local/bin/flex_inmem.py --live
WorkingDirectory=/usr/local/bin
Restart=on-failure
RestartSec=5
# 5 carriers x ~0.53x realtime ~= 2.65 cores; cap at 3 so this never eats the box.
CPUQuota=300%
# decoder threads live in the unit cgroup; killed with the flowgraph
KillMode=control-group

[Install]
WantedBy=multi-user.target
